# debian/sample.user.service
[Unit]
Description=sample

[Service]
ExecStart=/usr/bin/demo
Restart=always

[Install]
WantedBy= default.target